# etc/cron.d/ws2812XmlServerRestart
2 0 * * * root /opt/ws2812XmlServer/ws2812XmlServerRestart.sh
